# infra/kis.service
[Unit]
Description=KIS Server
After=network.target

[Service]
WorkingDirectory=/home/kis/kis
ExecStart=/home/kis/kis/kis -port 3000 -mode release
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
Type=simple
Restart=on-failure
User=kis
Group=www-data

[Install]
WantedBy=multi-user.target
Alias=kis.service
